# beeflow-server.socket
[Unit]
Description=Beeflow socket

[Socket]
ListenStream=/run/beeflow-server.sock

[Install]
WantedBy=sockets.target